# crontab
* * * * * /usr/local/bin/php /var/www/laravel/artisan schedule:run > /proc/1/fd/1 2>&1